# config/celery_beat.conf
[program:celery_beat]
command=/home/tixon/telegram/tg_base/venv/bin/celery -A core.celery.celery_app beat --loglevel=info
directory=/home/tixon/telegram/tg_base
user=tixon
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/tixon/telegram/tg_bot/logs/celery_beat_debug.log
stderr_logfile=/home/tixon/telegram/tg_bot/logs/celery_beat_error.log